# etc/cron.d/elk
0 * * * * root /opt/elk/bin/check_status.sh
