# etc/cron.d/elk
0 * * * * 		root /opt/elk/bin/check_status.sh
*/5 * * * *		root bash /opt/elk/bin/update_from_git.sh &> /var/log/git.log
*/10 * * * *	root > /log/logstash.log
*/10 * * * *	root > /log/logstash_2.log
*/10 * * * *	root > /log/logstash_3.log
*/10 * * * *	root > /log/logstash_4.log