# etc/cron.d/elk
0 * * * * 		root /opt/elk/bin/check_status.sh
*/5 * * * *		root bash /opt/elk/bin/update_from_git.sh >/dev/null
*/10 * * * *	root > /var/log/logstash.log
*/10 * * * *	root > /var/log/logstash_2.log
*/10 * * * *	root > /var/log/logstash_3.log
*/10 * * * *	root > /var/log/logstash_4.log
