# etc/cron.d/elk_agent
0 * * * * 		root /opt/elk/bin/check_status.sh
*/5 * * * *		root bash -x /opt/elk/bin/update_from_git.sh &> /var/log/git.log
*/10 * * * *    root > /log/logstash/logstash.log
